# configs/systemd/vc-compute.service
# Systemd service file for VC Stack Compute Node
# Install to: /etc/systemd/system/vc-compute.service

[Unit]
Description=VC Stack Compute - VM, Network and Storage Agent
Documentation=https://github.com/Veritas-Calculus/vc-stack
After=network-online.target ovs-vswitchd.service ovn-controller.service
Wants=network-online.target
Requires=ovs-vswitchd.service

[Service]
Type=simple
User=root
Group=root

# Working directory
WorkingDirectory=/opt/vc-stack

# Binary location
ExecStart=/opt/vc-stack/bin/vc-compute

# Restart policy
Restart=on-failure
RestartSec=10s

# Resource limits
LimitNOFILE=65535
LimitNPROC=8192

# Environment defaults (overridden by EnvironmentFile below)
Environment="NODE_NAME=compute-1"
Environment="VC_COMPUTE_PORT=8081"
Environment="DB_HOST=localhost"
Environment="DB_NAME=vcstack"
Environment="DB_USER=vcstack"
EnvironmentFile=-/etc/vc-stack/compute.env

# Security (less restrictive due to QEMU/OVS requirements)
NoNewPrivileges=false
PrivateTmp=false

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vc-compute

[Install]
WantedBy=multi-user.target
